# recipes-bsp/sensors/files/sensors-daemon.service
[Unit]
Description= Sensors Daemon Service
SourcePath=/usr/bin/sensors.qcom

[Service]
Restart=always
ExecStart=/usr/bin/sensors.qcom

[Install]
WantedBy=multi-user.target
